# scripts/systemd/jamessu-site-deploy.service
[Unit]
Description=Sync Notion + rebuild Fuwari site
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/jamessu_site
ExecStart=/srv/jamessu_site/scripts/deploy.sh
StandardOutput=journal
StandardError=journal
